# recipes-core/images/script-install/80-usbdisk.rules
KERNEL=="sd*[0-9]|sr*", SUBSYSTEM=="block", RUN+="/bin/usbdevinserted", SYMLINK+="usbdisk"
